# docker/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh

exec php-fpm